# ghost-browser/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1920x1080x24 -ac
autorestart=true
priority=10

[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -shared -rfbport 5900 -passwd %(ENV_VNC_PASSWORD)s -noxdamage -noxfixes
autorestart=true
priority=20
startsecs=2

[program:websockify]
command=/usr/bin/websockify --web=/usr/share/novnc 6080 localhost:5900
autorestart=true
priority=30
startsecs=3

[program:ghost-browser]
command=node /app/dist/index.js
directory=/app
user=ghost
autorestart=true
priority=40
startsecs=5
environment=PORT="%(ENV_PORT)s",DATA_DIR="%(ENV_DATA_DIR)s",API_KEY="%(ENV_API_KEY)s",DISPLAY=":99",HEADLESS="false",HOME="/home/ghost"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
